# distro/nova-updater/nova-server.service
[Unit]
Description=NOVA OS Web Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=nova
Group=nova
WorkingDirectory=/opt/nova-os
ExecStart=/usr/bin/node /opt/nova-os/server/index.js
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
